# internal/ruleengine/default.yaml
# KernelWatch default ruleset (embedded).
#
# This is the YAML translation of the original hardcoded Go detections. Operators
# extend/override it via KW_RULES_FILE / KW_RULES_DIR (append/override semantics);
# they never need to edit this file or recompile.
#
# Condition DSL (Falco-style): and/or/not, comparisons (= != & in contains
# startswith), $list references, and helpers (in_list, chain_in_list, has_prefix,
# contains_any, is_trusted_writer, is_docker_client, reverse_shell). The optional
# per-rule `lineage:` matrix escalates or suppresses by process ancestry.

lists:
  shells:                    [sh, bash, zsh, fish, dash, ash, ksh, tcsh]
  net_tools:                 [nmap, masscan, netcat, nc, ncat, socat, tcpdump, wireshark, tshark, curl, wget]
  high_net_tools:            [nmap, masscan]
  pkg_managers:              [apt, apt-get, dpkg, yum, dnf, rpm, apk, pip, pip3, npm, yarn, gem]
  priv_tools:                [sudo, su, nsenter, unshare, chroot, capsh, setuid, newgrp]
  interactive_benign_priv:   [sudo, su, newgrp]
  user_mgmt_tools:           [useradd, usermod, userdel, groupadd, groupmod, groupdel, chpasswd, passwd, adduser, deluser, gpasswd]
  log_tamper_tools:          [shred, wipe, srm]
  log_history_files:         [wtmp, btmp, lastlog, utmp, auth.log, secure, audit.log, .bash_history, .zsh_history, _history]
  host_tmp_dirs:             ["/tmp/", "/dev/shm/", "/var/tmp/"]
  writable_exec_dirs:        ["/tmp/", "/dev/shm/", "/var/tmp/", "/run/"]
  build_lineage:             [dpkg, apt, apt-get, rpm, yum, dnf, make, gcc, cc, go, cargo, npm, yarn, pip, pip3, cloud-init]
  sensitive_files:           ["/etc/shadow", "/etc/passwd", "/etc/sudoers", "/root/.ssh", "/var/run/docker.sock", "/.dockerenv", "/proc/sysrq-trigger", "/proc/kcore"]
  cred_files:                ["/.env", "/.aws/credentials", "/.gcloud/credentials", "/run/secrets", "/.kube/config"]
  persistence_files:         ["/etc/cron", "/var/spool/cron", "/etc/systemd/", "/lib/systemd/system/", "/usr/lib/systemd/system/", "/etc/init.d/", "/etc/rc.local", "/etc/update-motd.d/", "/etc/profile.d/", "/etc/ld.so.preload", "/etc/ld.so.conf", "/root/.ssh/authorized_keys", "/etc/sudoers.d/", "/etc/passwd", "/etc/shadow"]
  host_persistence_prefixes: ["/etc/cron", "/var/spool/cron", "/etc/systemd/", "/lib/systemd/system/", "/usr/lib/systemd/system/", "/etc/init.d/", "/etc/rc.local", "/etc/update-motd.d/", "/etc/profile.d/", "/etc/ld.so.preload", "/etc/ld.so.conf", "/etc/ssh/", "/etc/sudoers", "/etc/pam.d/", "/etc/modules-load.d/", "/root/.ssh/authorized_keys"]
  host_persistence_substrings: ["/.ssh/authorized_keys"]

macros:
  is_execve:    "evt.type = execve"
  spawns_shell: "is_execve and in_list(proc.exe_base, $shells)"

rules:
  - id: kernel_module_load
    desc: "Kernel module loaded"
    scope: all
    tactic: Persistence
    technique: T1547.006
    tags: [rootkit, kernel]
    condition: "evt.type = init_module"
    lineage:
      - when: any
        and: "scope = host and (is_trusted_writer or lineage = trusted)"
        action: suppress
      - when: any
        and: "scope = host"
        severity: high
        reason: "kernel module loaded on host from untrusted lineage (possible rootkit)"
      - when: any
        severity: critical
        reason: "kernel module loaded from container (possible rootkit)"

  - id: bpf_prog_load
    desc: "eBPF program loaded"
    scope: all
    tactic: Defense Evasion
    technique: T1562.001
    tags: [ebpf, evasion]
    condition: "evt.type = bpf and bpf.cmd = 5"
    lineage:
      - when: any
        and: "scope = host and is_trusted_writer"
        action: suppress
      - when: any
        severity: high
        reason: "eBPF program loaded (possible evasion/rootkit)"
        details: { bpf_cmd: bpf.cmd }

  - id: reverse_shell
    desc: "Reverse-shell / remote-exec pattern in command line"
    scope: all
    tactic: Execution
    technique: T1059
    tags: [reverse-shell, c2]
    condition: >-
      is_execve and proc.cmdline != "" and reverse_shell(proc.cmdline)
    severity: critical
    reason: "reverse-shell pattern detected in command line"
    details: { cmdline: proc.cmdline }

  - id: process_injection
    desc: "ptrace-based process injection/manipulation"
    scope: all
    tactic: Privilege Escalation
    technique: T1055.008
    tags: [injection]
    condition: "evt.type = ptrace and ptrace.request in [4, 5, 6, 13, 16, 0x4206]"
    severity: high
    reason: "ptrace process injection/manipulation"
    details: { request: ptrace.request, target_pid: ptrace.target }

  - id: persistence
    desc: "Write to a persistence-sensitive location"
    scope: all
    tactic: Persistence
    technique: T1543
    tags: [persistence]
    condition: >-
      evt.type = open and evt.is_open_write and (
        (scope = container and has_prefix(fd.name, $persistence_files)) or
        (scope = host and (has_prefix(fd.name, $host_persistence_prefixes) or
                           contains_any(fd.name, $host_persistence_substrings)))
      )
    lineage:
      - when: any
        and: "scope = host and is_trusted_writer"
        action: suppress
      - when: any
        and: 'fd.name startswith "/etc/ld.so.preload"'
        severity: critical
        reason: "write to persistence-sensitive location"
        details: { file: fd.name }
      - when: any
        severity: high
        reason: "write to persistence-sensitive location"
        details: { file: fd.name }

  - id: shell_in_container
    desc: "Shell executed inside a container"
    scope: container
    tactic: Execution
    technique: T1059
    tags: [container, shell]
    condition: spawns_shell
    lineage:
      - when: network
        severity: critical
        reason: "interactive shell spawned by network-facing service (possible RCE/web-shell)"
        details: { shell: proc.exe_base, spawned_by: lineage.network_parent }
      - when: trusted
        action: suppress
      - when: [unknown, interactive]
        severity: high
        reason: "shell execution inside container"
        details: { shell: proc.exe_base }

  - id: privilege_escalation
    desc: "Privilege-escalation tool executed"
    scope: all
    tactic: Privilege Escalation
    technique: T1548
    tags: [privesc]
    condition: "is_execve and in_list(proc.exe_base, $priv_tools)"
    lineage:
      - when: network
        severity: critical
        reason: "privilege escalation tool spawned by network-facing service"
        details: { tool: proc.exe_base, spawned_by: lineage.network_parent }
      - when: trusted
        action: suppress
      - when: interactive
        and: "in_list(proc.exe_base, $interactive_benign_priv)"
        action: suppress
      - when: interactive
        severity: high
        reason: "privilege/namespace tool executed in interactive session"
        details: { tool: proc.exe_base }
      - when: unknown
        severity: high
        reason: "privilege escalation tool executed"
        details: { tool: proc.exe_base }

  - id: network_tool
    desc: "Network/recon tool executed"
    scope: all
    tactic: Discovery
    technique: T1046
    tags: [network]
    condition: "is_execve and in_list(proc.exe_base, $net_tools)"
    lineage:
      - when: network
        severity: high
        reason: "network tool spawned by network-facing service (possible post-exploitation)"
        details: { tool: proc.exe_base, spawned_by: lineage.network_parent }
      - when: trusted
        action: suppress
      - when: interactive
        and: "in_list(proc.exe_base, $high_net_tools)"
        severity: high
        reason: "recon scanner executed in interactive session"
        details: { tool: proc.exe_base }
      - when: interactive
        action: suppress
      - when: unknown
        and: "in_list(proc.exe_base, $high_net_tools)"
        severity: high
        reason: "network tool executed"
        details: { tool: proc.exe_base }
      - when: unknown
        severity: medium
        reason: "network tool executed"
        details: { tool: proc.exe_base }

  - id: container_drift
    desc: "Execution of a binary from a writable/ephemeral path"
    scope: container
    tactic: Defense Evasion
    technique: T1036
    tags: [container, drift]
    condition: >-
      is_execve and proc.exe_path != "" and has_prefix(proc.exe_path, $writable_exec_dirs)
    severity: high
    reason: "binary executed from writable/ephemeral path"
    details: { path: proc.exe_path }

  - id: package_manager
    desc: "Package manager executed in a running environment"
    scope: all
    tactic: Execution
    technique: T1072
    tags: [package-manager]
    condition: "is_execve and in_list(proc.exe_base, $pkg_managers)"
    lineage:
      - when: network
        and: "scope = host"
        severity: critical
        reason: "package manager spawned by network-facing service (attacker installing tools)"
        details: { package_manager: proc.exe_base, spawned_by: lineage.network_parent }
      - when: network
        severity: high
        reason: "package manager spawned by network-facing service (attacker installing tools)"
        details: { package_manager: proc.exe_base, spawned_by: lineage.network_parent }
      - when: any
        and: "scope = host"
        action: suppress
      - when: trusted
        action: suppress
      - when: [unknown, interactive]
        severity: medium
        reason: "package manager executed inside running container"
        details: { package_manager: proc.exe_base }

  - id: sensitive_file
    desc: "Sensitive file accessed"
    scope: all
    tactic: Collection
    technique: T1005
    tags: [file]
    condition: >-
      evt.type = open and has_prefix(fd.name, $sensitive_files) and (
        scope = container or
        (evt.is_open_write and not fd.name = "/var/run/docker.sock"
         and not fd.name startswith "/root/.ssh")
      )
    lineage:
      - when: any
        and: 'fd.name = "/var/run/docker.sock"'
        severity: critical
        reason: "docker socket accessed by container (container-escape vector)"
        details: { file: fd.name }
      - when: any
        severity: medium
        reason: "sensitive file accessed"
        details: { file: fd.name }

  - id: credential_file
    desc: "Credential file accessed"
    scope: all
    tactic: Credential Access
    technique: T1552
    tags: [credentials]
    condition: "evt.type = open and contains_any(fd.name, $cred_files)"
    severity: high
    reason: "credential file accessed"
    details: { file: fd.name }

  - id: host_user_manipulation
    desc: "User/group account manipulation on host"
    scope: host
    tactic: Persistence
    technique: T1136.001
    tags: [host, accounts]
    condition: "is_execve and in_list(proc.exe_base, $user_mgmt_tools)"
    lineage:
      - when: network
        severity: critical
        reason: "account manipulation by network-facing service (backdoor account)"
        details: { tool: proc.exe_base, spawned_by: lineage.network_parent }
      - when: interactive
        severity: low
        reason: "account management in interactive admin session"
        details: { tool: proc.exe_base }
      - when: [unknown, trusted]
        severity: medium
        reason: "account manipulation on host"
        details: { tool: proc.exe_base }

  - id: host_log_tampering
    desc: "Log / forensic artifact tampering on host"
    scope: host
    tactic: Defense Evasion
    technique: T1070
    tags: [host, anti-forensics]
    condition: >-
      (is_execve and in_list(proc.exe_base, $log_tamper_tools)) or
      (evt.type = open and evt.is_open_write and
       (fd.name startswith "/var/log" or contains_any(fd.name, $log_history_files)))
    lineage:
      - when: any
        and: "is_execve and contains_any(proc.cmdline_lc, $log_history_files)"
        severity: high
        reason: "forensic artifact wiped with shred/wipe"
        details: { tool: proc.exe_base, cmdline: proc.cmdline }
      - when: any
        and: "is_execve"
        severity: medium
        reason: "file-wiping tool executed on host"
        details: { tool: proc.exe_base }
      - when: any
        and: "is_trusted_writer"
        action: suppress
      - when: any
        and: "evt.is_open_trunc or contains_any(fd.name, $log_history_files)"
        severity: high
        reason: "log/forensic file truncated or modified on host"
        details: { file: fd.name }
      - when: any
        severity: medium
        reason: "unexpected write to /var/log on host"
        details: { file: fd.name }

  - id: host_docker_sock
    desc: "Docker control socket opened by non-Docker process"
    scope: host
    tactic: Privilege Escalation
    technique: T1610
    tags: [host, docker, escape]
    condition: >-
      evt.type = open and (fd.name = "/var/run/docker.sock" or fd.name = "/run/docker.sock")
    lineage:
      - when: any
        and: "is_docker_client(proc.name)"
        action: suppress
      - when: any
        severity: high
        reason: "Docker control socket opened by non-Docker process on host (escape / lateral movement)"
        details: { process: proc.name, file: fd.name }

  - id: host_tmp_exec
    desc: "Execution from a world-writable temp dir on host"
    scope: host
    tactic: Execution
    technique: T1204.002
    tags: [host, drift]
    condition: >-
      is_execve and proc.exe_path != "" and has_prefix(proc.exe_path, $host_tmp_dirs)
      and not chain_in_list($build_lineage) and not is_trusted_writer
    severity: high
    reason: "binary executed from world-writable temp dir on host"
    details: { path: proc.exe_path }

  - id: host_shell_from_service
    desc: "Shell spawned by a network-facing service on host"
    scope: host
    tactic: Execution
    technique: T1059
    tags: [host, shell, rce]
    condition: spawns_shell
    lineage:
      - when: network
        severity: critical
        reason: "shell spawned by network-facing service on host (possible RCE)"
        details: { shell: proc.exe_base, spawned_by: lineage.network_parent }
